# infra/recipes/files/etc/cron.d/disable-turbo
@reboot root sh -c "echo 0 > /sys/devices/system/cpu/cpufreq/boost"
